# deploy/gunicorn.service
[Unit]
Description=RakshaGIS Gunicorn application server
Documentation=https://docs.gunicorn.org
After=network.target postgresql.service redis.service

[Service]
Type=notify
User=raksha
Group=raksha
WorkingDirectory=/opt/rakshagis/RakshaGIS
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
EnvironmentFile=/opt/rakshagis/RakshaGIS/.env
ExecStart=/opt/rakshagis/RakshaGIS/.venv/bin/gunicorn \
    --workers 4 \
    --bind unix:/run/rakshagis/gunicorn.sock \
    --worker-class sync \
    --timeout 120 \
    --access-logfile /opt/rakshagis/logs/gunicorn-access.log \
    --error-logfile /opt/rakshagis/logs/gunicorn-error.log \
    --log-level info \
    config.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
RuntimeDirectory=rakshagis
RuntimeDirectoryMode=0755
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
